# .config/autostart/gnome-terminal.desktop
[Desktop Entry]
Exec=gnome-terminal --maximize
Type=Application
Name[ja_JP]=たーみなる
Name=たーみなる
